# deploy_script/aas_edge_client.socket
[Unit]
Description=AAS Edge Client socket

[Socket]
ListenStream=/run/aas_edge_client.sock

[Install]
WantedBy=sockets.target